# challenges/corporate-network-assessment-lab/embedded-iot-device/etc/xinetd.d/telnet-custom
service telnet-custom
{
    disable     = no
    flags       = REUSE
    socket_type = stream
    wait        = no
    user        = root
    server      = /usr/sbin/in.telnetd
    server_args = -h -L /opt/corporate/bin/custom_shell.sh
    log_on_failure += USERID
    port        = 28
    bind        = 0.0.0.0
}